# project.conf
# 生产环境 supervisor 配置

[program:watersky]
command = gunicorn watersky.wsgi_test -b 127.0.0.1:8000
directory = /root/watersky/watersky
autostart = true
autorestart = true